# systemd/systemd-boot-btrfs.path
[Unit]
Description=Systemd-boot btrfs service monitoring /.snapshots/

[Path]
PathChanged=/.snapshots/
Unit=systemd-boot-btrfs.servive

[Install]
WantedBy=multi-user.target
